# backend/docker/supervisord.conf
; Slick API — supervisord orquesta los 4 procesos del contenedor.
; Todos loguean a stdout/stderr para que Docker/Dokploy los recoja.

; Sin "user=": el contenedor ya arranca como www-data (USER en el Dockerfile).
; Poner user= aquí haría que supervisord intente bajar privilegios siendo ya
; no-root y aborte con "can't drop privilege as nonroot user".
; logfile a stderr para que el orquestador (Dokploy) recoja sus diagnósticos.
[supervisord]
nodaemon=true
logfile=/dev/stderr
logfile_maxbytes=0
pidfile=/run/php/supervisord.pid

[program:php-fpm]
command=php-fpm --nodaemonize
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Worker de la cola Redis (confirmaciones/recordatorios por correo).
; --max-time=3600: el proceso se recicla cada hora (supervisord lo relanza)
; para recoger cambios de config y evitar fugas de memoria largas.
[program:queue-worker]
command=php /var/www/html/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
autorestart=true
priority=30
stopwaitsecs=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Scheduler de Laravel en primer plano (sustituye a cron): dispara
; appointments:send-reminders cada minuto (routes/console.php).
[program:scheduler]
command=php /var/www/html/artisan schedule:work
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
